# openrc/rc.local
#!bin/sh -e

chgrp power /etc/power/state
chmod g+w /etc/power/state

exit 0
